# RaspberryPiServer/MQTTFileScribe/MQTTFileScribe.service
[Unit]
Description=MQTT File Based Scribe
After=network.target
Requires=influxdb.service
Requires=mosquitto.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/local/skimon/bin/MQTTFileScribe.py

[Install]
WantedBy=multi-user.target